# dataserver/downloader/crontab
0 0 * * * curl http://subterfuge-game.com/leaderboards.html -o /data/leaderboard_$(date +%Y%m%d_%H%M%S).html